# oppf2-portal/trunk/oppf-portal/etc/crontab/daily_stockInfo_update.sh
#!/bin/sh

#dev
/usr/bin/curl http://127.0.0.1:8096/apt/stock/

#stage
#/usr/bin/curl http://127.0.0.1:8096/apt/stock/

#pord
#/usr/bin/curl http://127.0.0.1:8096/apt/stock/
